# config/gratia-gold.cron
15 * * * * root /usr/bin/gratia-gold -s 900
